# projects/VOLUMIO.krescue/OVERLAY.VDATA/dyn/etc/rc.local
#!/bin/bash

echo heartbeat > /sys/class/leds/sys_led/trigger

## hyphop ##
## fix cpu freq
## fix audio lags clicks

C=/sys/devices/system/cpu/cpufreq
[ "$FREQ" ] || FREQ=1200000

echo "[i] set fix freq $FREQ">&2

for c in $C/policy*; do
    echo $FREQ > $c/scaling_max_freq
    echo $FREQ > $c/scaling_min_freq
done

# /opt/cpu_freq

/usr/bin/amixer sset 'Audio hdmi-out mute' off
/usr/bin/amixer sset 'Audio spdif mute' off

exit 0
