# projects/VOLUMIO.krescue/OVERLAY.VDATA/dyn/etc/systemd/system/systemd-halt.service
[Unit]
Description=vim_poweroff
DefaultDependencies=no
After=final.target

[Service]
Type=oneshot
ExecStart=/opt/poweroff

[Install]
WantedBy=final.target
